# infrastructure/configs/systemd/dahouselab-backup.service
# daHouseLab nightly backup — installed to /etc/systemd/system/ by the
# execute-backup runbook. Runs the scripted runbook; on SUCCESS pings the
# Uptime Kuma push monitor (dead-man switch: no ping → Telegram alert).
[Unit]
Description=daHouseLab platform backup (docs/runbooks/execute-backup.md)
Wants=mnt-backups.mount
After=mnt-backups.mount docker.service

[Service]
Type=oneshot
ExecStart=/opt/dahouselab/scripts/backup/run-backup.sh
# Optional push URL (secret) — file is host-only, chmod 600, never in Git:
#   /etc/dahouselab/backup-push.env  ->  KUMA_PUSH_URL=https://status.<domain>/api/push/<token>
EnvironmentFile=-/etc/dahouselab/backup-push.env
# Runs ONLY if ExecStart succeeded — that is the dead-man contract
ExecStartPost=/bin/sh -c '[ -z "${KUMA_PUSH_URL:-}" ] || curl -fsS -m 10 "$KUMA_PUSH_URL" >/dev/null'
